# etc/systemd/system/wiki1c.service
[Unit]
Description=Read-only access to 1C knowledge base
After=network.target

[Service]
ExecStart=node --openssl-legacy-provider .
WorkingDirectory=/home/wiki1c/wiki1c
User=wiki1c
Group=wiki1c

[Install]
WantedBy=multi-user.target
